# drivemond-admin-new-install-3.1/deploy/supervisor/vito-worker.conf
; Supervisor program for the Vito queue worker.
; Required in production: the ride-timeout auto-cancel job (RideTimeoutJob) is
; queued and silently dropped unless a worker is running with QUEUE_CONNECTION=redis.
;
; Install:
;   sudo cp deploy/supervisor/vito-worker.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start vito-worker:*
;
; Adjust the paths/user below to your deployment.

[program:vito-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/vito/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/vito-worker.log
stopwaitsecs=3600
