# drivemond-admin-new-install-3.1/deploy/supervisor/vito-worker.conf
; Vito Laravel Queue Worker — Supervisor configuration
;
; Deploy steps:
;   1. Copy this file to /etc/supervisor/conf.d/vito-worker.conf
;   2. Adjust paths (VITO_PATH) and user (RUN_AS_USER) to match your setup
;   3. Run: sudo supervisorctl reread && sudo supervisorctl update
;   4. Start:  sudo supervisorctl start vito-worker:*
;
; The ride-timeout auto-cancel job needs this worker running — without it the job is
; silently dropped when QUEUE_CONNECTION=database/file/sync (it is never executed).
; For production use QUEUE_CONNECTION=redis and ensure REDIS_* vars are configured.
;
; See also: deploy/systemd/vito-worker.service (systemd alternative)

[program:vito-worker]
; Adjust to your deployment path
command=php /var/www/vito/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
; Run as the web-server user; never root in production.
user=www-data
autostart=true
autorestart=true
; Log stdout/stderr to Laravel's standard log channel.
stdout_logfile=/var/www/vito/storage/logs/supervisor-vito-worker.log
stderr_logfile=/var/www/vito/storage/logs/supervisor-vito-worker-error.log
; Restart immediately if the worker crashes.
startsecs=3
; Keep up to 3 instances — not needed for most Vito workloads.
numprocs=1
; Kill any worker still running after 1 hour to prevent memory leaks.
stopwaitsecs=30
stopasgroup=true
killasgroup=true
